# configs/supervisord.d/dci-lab.ini
[program:dci-lab]
command=/usr/bin/python3 /opt/dci-provisioner/api/lab.py
environment=LAB_HOST="%(ENV_LAB_HOST)s",LAB_PORT="%(ENV_LAB_PORT)s"
stdout_logfile=/var/log/supervisord/dci-lab.log
redirect_stderr=true
